# nvidia-driver-monitor-minimal.service
[Unit]
Description=NVIDIA Driver Package Status Web Server (Minimal Security)
Documentation=https://github.com/user/nvidia-driver-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nvidia-monitor
Group=nvidia-monitor
WorkingDirectory=/opt/nvidia-driver-monitor
ExecStart=/opt/nvidia-driver-monitor/nvidia-web-server -addr :8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nvidia-driver-monitor

# Environment variables for network timeouts
Environment=HTTP_TIMEOUT=60s
Environment=DIAL_TIMEOUT=30s
Environment=TLS_HANDSHAKE_TIMEOUT=30s

# Minimal security restrictions - prioritize network access
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=1024
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
